# livechain.ini
[program:livechain]
command=/usr/local/bin/python3.5 /opt/ebc/livechain/livechain.py 127.0.0.1 10515
redirect_stderr=true
autorestart=true
startretries=9999
